# jars-refresh.service
[Unit]
Description = Jars Refresh daemon

[Service]
Type = simple
ExecStart = /var/www/jars/plugins/jars-refreshd/jars-refresh.sh start
Restart = on-failure
User = www-data
WorkingDirectory = /var/www/jars/plugins/jars-refreshd

[Install]
WantedBy = multi-user.target
